# crontab

# 1H — every hour
0 * * * * /usr/local/bin/cronjob.sh 1h

# 4H — every 4 hours (00,04,08,12,16,20)
0 */4 * * * /usr/local/bin/cronjob.sh 4h

# 1D — every day at 19:00 UTC (3–5h before day end)
0 19 * * * /usr/local/bin/cronjob.sh 1d

# 1W — every Friday at 12:00 UTC (1–3 days before Monday)
0 12 * * 5 /usr/local/bin/cronjob.sh 1w

